# root/etc/services.d/octoprint/run
#!/usr/bin/with-contenv sh
exec octoprint --basedir /octoprint/octoprint config set --json webcam "{\"ffmpeg\":\"/usr/bin/ffmpeg\",\"snapshot\":\"http://10.1.2.101:8080/images/live.jpg\",\"stream\":\"http://10.1.2.101:8080/hls/live.stream.m3u8\"}" && octoprint serve --iknowwhatimdoing --host 0.0.0.0 --port 5000 --basedir /octoprint/octoprint 